# scripts/n8n-backup.cron
# n8n backup job - runs daily at 2 AM
0 2 * * * root /opt/apps/n8n/scripts/backup-n8n.sh